# .config/systemd/user/clean_downloads.service
[Unit]
Description=Run clean_downloads

[Service]
Type=oneshot
ExecStart=/home/arya/scripts/clean_downloads.sh
