# .config/systemd/user/clean_downloads.timer
[Unit]
Description=Clean the downloads folder daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
