# scripts/clean-download.timer
[Unit]
Description=Run clean-download daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
