# scripts/clean-download.service
[Unit]
Description=Clean Downloads Folder

[Service]
Type=oneshot
ExecStart=/home/arya/dotfiles/scripts/clean-download.sh
